# notify-battery.service
[Unit]
Description=System status notifications

[Service]
Type=oneshot
ExecStart=/opt/nudash/notify-battery
